# benchmark-build.service.template
# A benchmark target
#
# This shutsdown once the boot has completed

[Unit]
Description=Default
Requires=basic.target
After=basic.target
AllowIsolate=yes

[Service]
Type=oneshot
ExecStart=XXX_TEST_ROOT_XXX/build-dir.sh /root/stress-ng.git/
ExecStartPost=/sbin/poweroff

[Install]
WantedBy=multi-user.target
